# task3/cleanup.timer
[Unit]
Decription=Every 2 min log

[Timer]
Type=oneshot
OnCalendar=*-*-* *:0/2
Presistent=true

[Install]
WantedBy=timers.target








